# supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=php5-fpm

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND